# debian/conf/gunicorn_ddrlocal.conf
[program:ddrlocal]
user=ddr
directory=/usr/local/src/ddr-local/ddrlocal
command=gunicorn ddrlocal.wsgi:application
autostart=true
autorestart=true
redirect_stderr=True
